# etc/supervisor.d/tftp.ini
[program:tftp]
directory = /netboot/
command = /usr/local/bin/py3tftp --host 0.0.0.0 -p 69 
daemon=off
priority = 1
stdout_logfile=/scripts/tftpd.log
redirect_stderr=true
