# etc/supervisor.d/webserver.ini
[program:webserver]
command = /usr/bin/python /bin/webserver.py
startretries = 2
daemon = off
priority = 2
