# etc/supervisor.d/webserver.ini
[program:webserver]
command = /usr/local/bin/python /scripts/webserver.py
startretries = 2
daemon = off
priority = 2
