# data/etc/supervisor/supervisor.d/tftp.ini
[program:tftp]
directory = /tftp/netboot/
command = /usr/bin/py3tftp --host 0.0.0.0 -p 69 
daemon=off
priority = 1
stdout_logfile=/tmp/tftpd.log
redirect_stderr=true
