# data/etc/supervisor/supervisor.d/webserver.ini
[program:webserver]
command = /usr/bin/python3 /scripts/webserver.py
startretries = 2
daemon = off
priority = 2
